# packaging/nuc-powerd.service
[Unit]
Description=NUC Power Optimization Daemon
After=multi-user.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nuc-powerd
ExecStart=/usr/local/bin/nuc-powerd --config /etc/nuc-powerd.toml run
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ReadWritePaths=/run/nuc-powerd

[Install]
WantedBy=multi-user.target
